# rootfs_overlay/etc/inittab
# /etc/inittab
# https://unix.stackexchange.com/questions/299408/how-to-login-automatically-without-typing-root-in-buildroot-x86-64-qemu
console::respawn:/bin/sh
